# iptables/templates/default/iptables.service.erb
[Unit]
Description=Spark service

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=<%= @start_script %>

[Install]
WantedBy=multi-user.target

