# mongrel2/templates/default/mongrel2.conf.erb
start on ((local-filesystems and runlevel [2345]) or vagrant-mounted)
stop on runlevel [016]

<%#start on (started networking)
  stop on runlevel [!06]%>

console output

pre-start script
  <% @servers.each do |server_name| %>
  start m2-<%= server_name %>
  <% end %>
end script

