# postgres/templates/default/pgbouncer/upstart.conf.erb
description "pgbouncer upstart script"
start on (local-filesystems and runlevel [2345])
stop on runlevel [016]
respawn
limit nofile 65000 65000
exec <%= @cmd %>

