# spiped/templates/default/upstart.conf.erb
###############################################################################
# spiped daemon script
#
# spiped allows secure communication between 2 boxes that don't have ssl
###############################################################################
start on (local-filesystems and runlevel [2345])
stop on runlevel [016]

pre-start script
    mkdir -p -m0777 <%= @run_dir %>
    chown <%= @username %>:<%= @group %> <%= @run_dir %>
end script

respawn
expect fork
exec start-stop-daemon --start --chuid <%= @username %>:<%= @usergroup %> --pidfile <%= @pidfile %> --exec <%= @exec %> -- <%= @args %>

