# spiped/templates/default/upstart.conf.erb
# spiped daemon script
#
# spiped allows secure communication between 2 boxes that don't have ssl

start on (local-filesystems and runlevel [2345])
stop on runlevel [016]

<% if defined?(@username) && !@username.empty? %>
# setuid <%= @username %>
<% end %>
<% if defined?(@usergroup) && !@usergroup.empty? %>
# setgid <%= @usergroup %>
<% end %>

pre-start script
    # mode is world writeable because the pid file needs to be written to the spiped dir
    # mode is world executable because evidently you need to execute something to write
    # directory is completely opened because each spiped command could be run under a different user
    mkdir -p -m0777 <%= @run_dir %>
    #chown <%= @username %>:<%= @group %> <%= @run_dir %>
end script

#expect fork
#exec <%= @exec %> <%= @args %>

respawn
expect fork
exec start-stop-daemon --start --chuid <%= @username %>:<%= @usergroup %> --pidfile <%= @pidfile %> --exec <%= @exec %> -- <%= @args %>

