# spiped/templates/default/upstarts.conf.erb
description "spiped catch-all, this will start all registered spiped scripts at once"
start on ((local-filesystems and runlevel [2345]) or vagrant-mounted)

pre-start script
  <% @services.each do |service_name| %>
    if [ -n "$(status <%= service_name %> | grep 'stop/waiting')" ]; then
        initctl start <%= service_name %>
    fi
  <% end %>
end script

post-stop script
  <% @services.each do |service_name| %>
    if [ -n "$(status <%= service_name %> | grep 'start/running')" ]; then
        initctl stop <%= service_name %>
    fi
  <% end %>
end script

