# daemon/templates/default/instances.conf.erb
stop on runlevel [06]
stop on stopping <%= @recipe_name %>
start on starting <%= @recipe_name %>

env NUM_WORKERS=<%= @count %>

pre-start script
  for i in $(seq 1 $NUM_WORKERS); do
    start child-<%= @instance_name %> N=$i
  done
end script

