# daemon/templates/default/service@.service.erb
[Unit]
Description="Daemon <%= @service_name %> instance #%i"
PartOf=<%= @partof %>

[Service]
Type=simple
Restart=on-success
RestartSec=5
LimitNOFILE=100000
<% if defined?(@chdir) %>
WorkingDirectory=<%= @chdir %>
<% end %>
<% @environ_vars.each do |environ| %>
Environment="<%= environ %>"
<% end %>
<% @environ_files.each do |environ| %>
EnvironmentFile=<%= environ %>
<% end %>
<% if defined?(@username) %>
User=<%= @username %>
<% end %>
<% if defined?(@group) %>
Group=<%= @group %>
<% end %>
ExecStart=<%= @command %>

