# sentry/templates/default/service.conf.erb
start on started network-services
stop on stopping network-services
respawn

setuid <%= @user %>
setgid <%= @user %>

script
  exec <%= @command %>
end script

