# redis/templates/default/redis-server.service.erb
[Unit]
Description="Redis systemd script for redis chef cookbook"
After=syslog.target network.target remote-fs.target
Requires=network-online.target

[Service]
Type=notify
NotifyAccess=all
Restart=on-success
RestartSec=5
LimitNOFILE=65000
<% if defined?(@username) %>
User=<%= @username %>
<% end %>
<% if defined?(@group) %>
Group=<%= @group %>
<% end %>
ExecStart=<%= @command %> <%= @config %>
ExecStop=<%= @command_shutdown %> shutdown

[Install]
WantedBy=mutli-user.target

