# uwsgi/templates/default/servers.conf.erb
description "uwsgi catch-all, this will start all registered uwsgi scripts at once"
start on (local-filesystems and runlevel [2345])

pre-start script
  <% @server_names.each do |server_name| %>
    if [ -n "$(status <%= server_name %> | grep 'stop/waiting')" ]; then
        initctl start <%= server_name %>
    fi
  <% end %>
end script

post-stop script
  <% @server_names.each do |server_name| %>
    if [ -n "$(status <%= server_name %> | grep 'start/running')" ]; then
        initctl stop <%= server_name %>
    fi
  <% end %>
end script
