# uwsgi/templates/default/server.conf.erb
# simple uWSGI script
# http://uwsgi-docs.readthedocs.org/en/latest/Upstart.html

description "uwsgi upstart for <%= @server_name %>"
start on ((local-filesystems and runlevel [2345]) or vagrant-mounted)
start on starting uwsgi
stop on runlevel [06]
stop on stopping uwsgi

<% if defined?(@chdir) %>
chdir <%= @chdir %>
<% end %>

respawn
respawn limit 5 5
limit nofile 100000 100000

#console output

script
<% if defined?(@environ) %>
  <%= @environ %>
<% end %>
  exec <%= @exec_str %>
end script


