# uwsgi/templates/default/servers.conf.erb
description "uwsgi catch-all, this will start all registered uwsgi scripts at once"
start on ((local-filesystems and runlevel [2345]) or vagrant-mounted)

